# upstart_jobs/sbr_daemon.conf
description "Demonio de adquisición de datos para el HMI del SBR"
author "Joaquín Sorianello <joac@joac.com.ar>"
start on runlevel [345]
stop on runlevel [0126]
exec /home/hmi/ArakurWW/adquisitor.py
