# upstart_jobs/sbr_web.conf
description "Web de control del SBR"
author "Joaquín Sorianello <joac@joac.com.ar>"
start on runlevel [2345]
stop on runlevel [016]
exec /home/hmi/ArakurWW/web.py
